# mariadb.ini
[program:mariadb]
command=mysqld_safe --basedir=/usr
